# crontab
0 23,7,9,13 * * * /usr/local/bin/python /CROUStillant/__main__.py >> /CROUStillant/crontab.log
