# crontab
0 1,9,11,15 * * * /usr/local/bin/python /CROUStillant/__main__.py >> /CROUStillant/crontab.log
