# crontab
0 23,7,9,13 * * * uv run --project /CROUStillant /CROUStillant/__main__.py >> /CROUStillant/crontab.log